# src/features/dc-s6-service-sshd/s6-rc.d/sshd-logger/run
#!/command/execlineb -PW
foreground { echo "[sshd] Starting logger ..." }
sed --unbuffered "s/^/[sshd] /"
